# acmeshomeassistant/rootfs/etc/services.d/crond/run
#!/usr/bin/with-contenv sh

echo "INFO: Starting crond service for certificate renewal"
exec crond -f
